# systemd/asus-control.service
[Unit]
Description=Linux ASUS Control automation daemon
After=multi-user.target

[Service]
Type=simple
WorkingDirectory=/opt/asus-control
ExecStart=/opt/asus-control/.venv/bin/python -m asus_control.daemon --config /opt/asus-control/config/config.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
